# .crontab
0 0 * * * cd ~/apps/redmine && make dump > logs/cron.log 2>&1
